# systemd/insight-server.service
[Unit]
Description=in_sight kiosk server (Flask via waitress)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/opt/in_sight
# Secrets store (Anthropic key, dVerse password) — written by the Settings page.
# The leading "-" makes it optional (no error if it doesn't exist yet).
EnvironmentFile=-/opt/in_sight/config/.env
# Build a fresh snapshot on boot, then serve it.
ExecStartPre=/opt/in_sight/scripts/refresh.sh
ExecStart=/opt/in_sight/.venv/bin/waitress-serve --host=0.0.0.0 --port=8137 server.app:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
